# data/profiles/img/extra/etc/systemd/system/getty@tty1.service.d/autologin.conf
# could consider for safety to login only once we have enough disk space after growfs. this should not really be necessary tho.
# (note: growfs is only part of efly-img and would error for efly-dd)
#[Unit]
#After=boot-growfs.service

[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin efly - $TERM
