# data/img/extra/etc/systemd/system/reflector.service.d/require-growfs.conf
[Unit]
After=network-online.target nss-lookup.target boot-growfs.service
